# otmonitor/root/etc/s6-overlay/s6-rc.d/otmonitor/run
#!/command/with-contenv bashio
set -e

additional_options=$(bashio::config 'additional_options')

bashio::log.info 'Started otmonitor'
# shellcheck disable=SC2086
/usr/local/bin/otmonitor -a --daemon --system -w 8080 -p 7686 -f $additional_options
